# byos-terminus/rootfs/etc/s6-overlay/s6-rc.d/puma/run
#!/command/with-contenv bash
# shellcheck shell=bash
# Puma web server with inline database migrations

set -e

cd /app

# Load environment from bashio (HA config)
export HANAMI_ENV="${HANAMI_ENV:-production}"
export RACK_ENV="${RACK_ENV:-production}"
# Puma listens on 3000, nginx proxies from 2300
export HANAMI_PORT="3000"
export HANAMI_SERVE_ASSETS=true

# Construct DATABASE_URL for local PostgreSQL
export DATABASE_URL="postgres://postgres@localhost/terminus"

# Construct REDIS_URL for local Valkey (Unix socket for performance)
export REDIS_URL="unix:///var/run/valkey/valkey.sock"

# Also set KEYVALUE_URL (Terminus naming)
export KEYVALUE_URL="$REDIS_URL"

# Load APP_SECRET from generated file or generate
if [ -f /data/.app_secret ]; then
    export APP_SECRET=$(cat /data/.app_secret)
fi

# Load API_URI from HA config or use default
# Note: External port is 2300 (nginx), internal Puma is 3000
if [ -f /data/options.json ] && [ -n "$(jq -r '.api_uri // empty' /data/options.json 2>/dev/null)" ]; then
    export API_URI=$(jq -r '.api_uri' /data/options.json)
else
    # Default to external nginx port - user should configure in HA options for device connectivity
    export API_URI="http://localhost:2300"
fi

echo "[puma] Environment: HANAMI_ENV=$HANAMI_ENV, PORT=$HANAMI_PORT"

# Create database if it doesn't exist
echo "[puma] Ensuring database exists..."
PSQL="/usr/lib/postgresql/16/bin/psql"
$PSQL -U postgres -tc "SELECT 1 FROM pg_database WHERE datname = 'terminus'" | grep -q 1 || \
    $PSQL -U postgres -c "CREATE DATABASE terminus"

# Run database migrations (as app user for proper permissions)
echo "[puma] Running database migrations..."
s6-setuidgid app bundle exec hanami db migrate
echo "[puma] Migrations complete"

# Compile assets if not already compiled (first run only)
if [ ! -d "/app/public/assets/_authentication" ]; then
    echo "[puma] Compiling assets (first run)..."
    bundle exec hanami assets compile 2>/dev/null || echo "[puma] Warning: Asset compilation may have issues"
    echo "[puma] Assets compiled"
fi

echo "[puma] Starting Puma web server..."

# Run as app user
exec s6-setuidgid app \
    bundle exec puma --config ./config/puma.rb
