# byos-terminus/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/command/with-contenv bash
# shellcheck shell=bash
# Nginx graceful shutdown

echo "[nginx] Service stopping (exit code: $1, signal: $2)"
exit 0
