# byos-terminus/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bash
# shellcheck shell=bash
# Nginx reverse proxy for HA ingress support

set -e

NGINX_CONF="/etc/terminus/nginx.conf"
PUMA_PORT=3000
MAX_WAIT=60

echo "[nginx] Waiting for Puma to be ready on port $PUMA_PORT..."

# Wait for Puma to be listening
WAITED=0
while ! curl -sf "http://127.0.0.1:${PUMA_PORT}/up" > /dev/null 2>&1; do
    if [ $WAITED -ge $MAX_WAIT ]; then
        echo "[nginx] ERROR: Puma not ready after ${MAX_WAIT}s, starting anyway..."
        break
    fi
    sleep 1
    WAITED=$((WAITED + 1))
done

if [ $WAITED -lt $MAX_WAIT ]; then
    echo "[nginx] Puma is ready after ${WAITED}s"
fi

echo "[nginx] Starting nginx reverse proxy..."

# Test config before starting
nginx -t -c "$NGINX_CONF" 2>&1 || {
    echo "[nginx] Configuration error, check $NGINX_CONF"
    exit 1
}

# Run nginx in foreground
exec nginx -g "daemon off;" -c "$NGINX_CONF"
