# byos-terminus/rootfs/etc/s6-overlay/s6-rc.d/postgres/finish
#!/command/with-contenv bash
# shellcheck shell=bash
# PostgreSQL graceful shutdown

echo "[postgres] Service stopping (exit code: $1, signal: $2)"

# Request graceful shutdown via pg_ctl if postgres is running
if pgrep -x postgres > /dev/null 2>&1; then
    echo "[postgres] Requesting graceful shutdown..."
    s6-setuidgid postgres /usr/lib/postgresql/16/bin/pg_ctl stop -D /data/postgres -m fast -w -t 30 2>/dev/null || true
fi

exit 0
